# ops/systemd/agentern-restore-verify.timer
[Unit]
Description=Verify Agentern backup restoration monthly

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
