# ops/systemd/agentern-backup.timer
[Unit]
Description=Run Agentern PostgreSQL backups daily

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
